# curator4-1072903490.conf
[program:curator4-1072903490]
command=python3 /home/curator4/curator4-1072903490/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1072903490/test.err.log
stdout_logfile=/home/curator4/curator4-1072903490/test.out.log